# deploy/systemd/amazon-spapi-catalog-us.service
[Unit]
Description=Amazon SP-API — US marketplace catalog sync
Documentation=file:///opt/amazon-spapi-workers/docs/生产部署.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=amazon-spapi
Group=amazon-spapi
WorkingDirectory=/opt/amazon-spapi-workers
EnvironmentFile=/etc/amazon-spapi/amazon-spapi.env

ExecStart=/home/amazon-spapi/.local/bin/uv run celery -A amazon_spapi.worker worker \
  -l info \
  -c 1 \
  -Q SpapiCatalogItemsUpdate_US

Restart=always
RestartSec=15
KillSignal=SIGTERM
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
