# deploy/systemd/amazon-spapi-offers-us.service
[Unit]
Description=Amazon SP-API — US marketplace offer refresh
Documentation=file:///opt/amazon-spapi-workers/docs/生产部署.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=amazon-spapi
Group=amazon-spapi
WorkingDirectory=/opt/amazon-spapi-workers
EnvironmentFile=/etc/amazon-spapi/amazon-spapi.env

# 调整 uv 路径：which uv 或 ~/.local/bin/uv
ExecStart=/home/amazon-spapi/.local/bin/uv run celery -A amazon_spapi.worker worker \
  -l info \
  -c 2 \
  -Q SpapiItemOffersUpdate_US

Restart=always
RestartSec=15
KillSignal=SIGTERM
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
